# playbooks/roles/infinidb_supervisor/templates/supervisord.conf.j2
[supervisord]
nodaemon=true
user={{ infinidb_user }}

[program:infinidb]
command={{ infinidb_supervisor_script }}